# day12/deploy/systemd/sigma-schema-agent.service
[Unit]
Description=Sigma Intelligence Platform — Schema Evolution Agent
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/sigma-platform
EnvironmentFile=/home/ec2-user/sigma-platform/.env
ExecStart=/usr/bin/python3 -m uvicorn lab.agents.schema_agent_app:app --host 0.0.0.0 --port 8002 --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sigma-schema-agent

[Install]
WantedBy=multi-user.target
